# packaging/deb/DEBIAN/postinst
#!/bin/sh
set -e

echo ""
echo "To enable LaNotifica for your user, run:"
echo "  systemctl --user enable --now lanotifica"
echo ""

exit 0
